# deploy/docker/mcp/supervisor/conf.d/mcp-servers.conf
# MCP Servers Supervisor Configuration
# Combined configuration for all MCP servers

# MCP Demo Server
[program:mcp-demo-server]
command=bash -c "uv sync && uv run server.py"
directory=/export/App/code/demo-server
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONIOENCODING="utf-8"
user=root

# MCP Scanner Server
[program:mcp-scanner]
command=bash -c "uv sync && source .venv/bin/activate && uv run uvicorn codes.api.main:app --port 8002 --host 0.0.0.0"
directory=/export/App/code/mcp_scanner
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONIOENCODING="utf-8"
user=root

# MCP JEB Server
[program:mcp-jeb]
command=bash /export/App/code/jebmcp/start_mcp_jeb.sh
directory=/export/App/code/jebmcp
autostart=true
autorestart=true
startsecs=30
startretries=3
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONIOENCODING="utf-8",JEB_SERVER_PORT="16161",JEB_MCP_PORT="8008"
user=root
priority=999



